# launch/com.opencode.statusball.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    LaunchAgent template for OpenCodeStatusBall.

    DO NOT install this file directly. The placeholders __HOME__ and __BIN__
    are replaced at install time by ./launch/install.sh.

    Modern install/uninstall (macOS 10.11+):
        launchctl bootstrap gui/<uid> <plist>
        launchctl bootout    gui/<uid>/<label>
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opencode.statusball</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BIN__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/OpenCode/statusball.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/OpenCode/statusball.err.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.opencode.statusball</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
